# data/systemd/system/incus.socket.d/kapsule-socket-group.conf
# SPDX-FileCopyrightText: 2024-2026 KDE Community
# SPDX-License-Identifier: GPL-3.0-or-later
#
# Drop-in for incus.socket to use wheel group

[Socket]
SocketGroup=wheel
