# images/mkosi.images/kapsule-dev/mkosi.extra/etc/profile.d/local-bin.sh
# SPDX-FileCopyrightText: 2026 Lasath Fernando <devel@lasath.org>
#
# SPDX-License-Identifier: GPL-3.0-or-later

# Add ~/.local/bin to PATH for tools installed per-user (e.g. kde-builder)
if [ -d "$HOME/.local/bin" ]; then
    PATH="$HOME/.local/bin:$PATH"
fi
