# data/systemd/user/flatpak-session-helper.service.d/kapsule-host-command.conf
# SPDX-FileCopyrightText: 2026 Lasath Fernando <devel@lasath.org>
# SPDX-License-Identifier: GPL-3.0-or-later
#
# Drop-in for flatpak-session-helper.service to enable Kapsule integration.
#
# When FLATPAK_HOST_COMMAND_WRAPPER is set, the patched session-helper
# prepends the wrapper path before every HostCommand argv, allowing
# Kapsule to route commands into the appropriate container.

[Service]
Environment=FLATPAK_HOST_COMMAND_WRAPPER=/usr/lib/kapsule/kapsule-flatpak-wrapper
